# DEB/DEBIAN/postinst
#!/bin/bash

update-rc.d vncsrv defaults
update-rc.d vncsrv enable
/etc/init.d/vncsrv start

exit 0
